# conf/etc/conf.d/postgres.conf
[program:postgresql]
command=/usr/lib/postgresql/9.3/bin/postgres -D /var/lib/postgresql/9.3/main -c config_file=/etc/postgresql/9.3/main/postgresql.conf
user=postgres
autorestart=true

#directory=/usr/local/
#autostart=true
#redirect_stderr=true

stdout_logfile=/var/log/supervisor/postgresql-stdout.log
stderr_logfile=/var/log/supervisor/postgresql-stderr.log
