# conf/old/supervisord.conf/gunicorn.conf
[program:gunicorn]
command=/home/dev/gunicorn testproj.wsgi:application --bind 127.0.0.1:8000 --pid /tmp/gunicorn.pid ;
directory=/tmp/testproj/ ;